# config/trollduction/supervisord_fernsehbild.conf
; pytroll supervisor config file to be included in main config file.
; 
; defines processes to "Fernsehbild" products
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;


; --------------------------------------------------
; ---------- Meteosat 10 ---------------------------
; --------------------------------------------------

; trollstalker defined in main supervisord.conf

; for products depending on nwcsaf cloudtype
[program:l2processor-msg-hrit-with-ct]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %L2PROCESSORPATH% -c %MASTERCONFIGPATH% -C l2_msg_hrit_with_cloudtype -N 100    
process_name=trollduction-l2processor-msg-hrit-with-ct
priority=3                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log

