# config/trollduction/supervisord_worldcomp.conf
; pytroll supervisor config file to be included in main config file.
; 
; defines processes to create world composite products
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;


; --------------------------------------------------
; ---------- Meteosat 10 ---------------------------
; --------------------------------------------------

; trollstalker defined in main supervisord.conf

[program:l2processor-msg-hrit-wcm]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %L2PROCESSORPATH% -c %MASTERCONFIGPATH% -C l2_msg_hrit -N 200    
process_name=trollduction-l2processor-msg-hrit-wcm
priority=3                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log


; --------------------------------------------------
; ---------- GOES 13 -------------------------------
; --------------------------------------------------
[program:trollstalker-goes13]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %TROLLSTALKERPATH% -c %MASTERCONFIGPATH% -C stalker_goes13    
process_name=trollduction-trollstalker-goes13
priority=5                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log

[program:seggather-goes13]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %SEGMENTGATHERERPATH% -c %MASTERCONFIGPATH% -C seg_gatherer_goes13
process_name=trollduction-seggather-goes13
priority=7                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log

[program:l2processor-goes13]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %L2PROCESSORPATH% -c %MASTERCONFIGPATH% -C l2_goes13
process_name=trollduction-l2processor-goes13
priority=6                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log

; --------------------------------------------------
; ---------- GOES 15 -------------------------------
; --------------------------------------------------
[program:trollstalker-goes15]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %TROLLSTALKERPATH% -c %MASTERCONFIGPATH% -C stalker_goes15   
process_name=trollduction-trollstalker-goes15
priority=5                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log

[program:seggather-goes15]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %SEGMENTGATHERERPATH% -c %MASTERCONFIGPATH% -C seg_gatherer_goes15
process_name=trollduction-seggather-goes15
priority=7                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log

[program:l2processor-goes15]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %L2PROCESSORPATH% -c %MASTERCONFIGPATH% -C l2_goes15
process_name=trollduction-l2processor-goes15
priority=6                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log

; --------------------------------------------------
; ---------- Himawari 8-----------------------------
; --------------------------------------------------
[program:trollstalker-himawari8]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %TROLLSTALKERPATH% -c %MASTERCONFIGPATH% -C stalker_himawari8
process_name=trollduction-trollstalker-himawari8
priority=5                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log

[program:seggather-himawari8]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %SEGMENTGATHERERPATH% -c %MASTERCONFIGPATH% -C seg_gatherer_himawari8
process_name=trollduction-seggather-himawari8
priority=7                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log

[program:l2processor-himawari8]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %L2PROCESSORPATH% -c %MASTERCONFIGPATH% -C l2_himawari8
process_name=trollduction-l2processor-himawari8
priority=6                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log


; --------------------------------------------------
; ---------- MTP -----------------------------------
; --------------------------------------------------
[program:trollstalker-mtp]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %TROLLSTALKERPATH% -c %MASTERCONFIGPATH% -C stalker_mtp
process_name=trollduction-trollstalker-mtp
priority=5                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log

[program:seggather-mtp]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %SEGMENTGATHERERPATH% -c %MASTERCONFIGPATH% -C seg_gatherer_mtp
process_name=trollduction-seggather-mtp
priority=7                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log

[program:l2processor-mtp]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %L2PROCESSORPATH% -c %MASTERCONFIGPATH% -C l2_mtp
process_name=trollduction-l2processor-mtp
priority=6                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log

; --------------------------------------------------
; ---------- Common --------------------------------
; --------------------------------------------------

[program:seggather-worldcomp-wcm3km-ir]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %SEGMENTGATHERERPATH% -c %MASTERCONFIGPATH% -C seg_gatherer_world_comp_wcm3km_ir
process_name=trollduction-seggather-worldcomp-wcm3km-ir
priority=7                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log

[program:seggather-worldcomp-wwwWcm40km-ir]
environment = PPP_CONFIG_DIR="%CONFIGDIR%"
command=%PYTHONBIN% %SEGMENTGATHERERPATH% -c %MASTERCONFIGPATH% -C seg_gatherer_world_comp_wwwWcm40km_ir
process_name=trollduction-seggather-worldcomp-wwwWcm40km-ir
priority=7                  
startsecs=3
stdout_logfile=%LOGDIR%/%(program_name)s-stdout.log
stderr_logfile=%LOGDIR%/%(program_name)s-stderr.log
